# scripts/logger.service
[Unit]
Description=MFC Logging Utility

[Service]
Type=simple
ExecStart=/usr/bin/logger
Restart=always

[Install]
WantedBy=multi-user.target
